# root/etc/services.d/suitarr/run
#!/usr/bin/with-contenv bash

source /etc/env

if [[ $APP == "radarr" ]]; then
    exec s6-setuidgid hotio mono --debug /app/Radarr.exe --nobrowser --data=/config $ARGS
fi

if [[ $APP == "sonarr" ]]; then
    exec s6-setuidgid hotio mono --debug /app/NzbDrone.exe --nobrowser --data=/config $ARGS
fi

if [[ $APP == "jackett" ]]; then
    exec s6-setuidgid hotio mono --debug /app/JackettConsole.exe --NoRestart --ListenPublic --Port 8080 --DataFolder=/config $ARGS
fi

if [[ $APP == "nzbhydra" ]]; then
    exec s6-setuidgid hotio python /app/nzbhydra.py --nobrowser --port 8080 --config /config/nzbhydra.cfg --database /config/nzbhydra.db --logfile /config/nzbhydra.log $ARGS
fi

if [[ $APP == "nzbget" ]]; then
    exec s6-setuidgid hotio /app/nzbget --server --option OutputMode=log --option ControlPort=8080 --configfile /config/nzbget.conf $ARGS
fi
